# ops/sub2api-watchdog.timer
[Unit]
Description=Run Sub2API watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=sub2api-watchdog.service

[Install]
WantedBy=timers.target
